# keyd.service
[Unit]
Description=key remapping daemon

[Service]
Type=simple
ExecStart=/usr/bin/keyd

[Install]
WantedBy=multi-user.target

